# systemd/grape-daily-upload.service
[Unit]
Description=GRAPE Daily DRF Upload to PSWS
Documentation=https://github.com/wsprdaemon/signal-recorder
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=wsprdaemon
Group=wsprdaemon

# Script location
ExecStart=/home/wsprdaemon/signal-recorder/scripts/daily-drf-upload.sh

# Environment - switch to production when ready
Environment=DATA_ROOT=/tmp/grape-test
Environment=CONFIG_FILE=/home/wsprdaemon/signal-recorder/config/grape-config.toml
Environment=VENV_PATH=/home/wsprdaemon/signal-recorder/venv
Environment=LOG_FILE=/tmp/grape-test/logs/daily-upload.log

# Production settings (uncomment when ready):
# Environment=DATA_ROOT=/var/lib/signal-recorder

# Timeouts - allow plenty of time for upload
TimeoutStartSec=3600
TimeoutStopSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grape-daily-upload

[Install]
WantedBy=multi-user.target
