# systemd/grape-daily-upload.service
[Unit]
Description=GRAPE Daily DRF Upload to PSWS
Documentation=https://github.com/mijahauan/grape-recorder
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

# User/Group - same as other grape services
User=grape-recorder
Group=grape-recorder

# Environment from file (set by install.sh)
EnvironmentFile=/etc/grape-recorder/environment

# Daily upload script
ExecStart=${GRAPE_PROJECT}/scripts/daily-drf-upload.sh

# Timeouts - allow plenty of time for upload
TimeoutStartSec=3600
TimeoutStopSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grape-daily-upload

[Install]
WantedBy=multi-user.target
