# systemd/grape-spectrograms.service
[Unit]
Description=GRAPE Spectrogram Generation - Daily with Solar Zenith
Documentation=https://github.com/mijahauan/grape-recorder
After=network.target grape-analytics.service

[Service]
Type=oneshot

# User/Group - same as other grape services
User=wsprdaemon
Group=wsprdaemon

# Environment from file
# EnvironmentFile=/etc/grape-recorder/environment

WorkingDirectory=/home/wsprdaemon/grape-recorder

# Generate daily spectrograms for today (with solar zenith + power graph)
# Runs every 10 minutes to update as data comes in throughout the day
ExecStart=/opt/grape-recorder/venv/bin/python -m grape_recorder.grape.carrier_spectrogram \
    --data-root /var/lib/grape-recorder \
    --all-channels \
    --date today \
    --grid EM38ww

StandardOutput=journal
StandardError=journal
SyslogIdentifier=grape-spectrograms

# Resource limits - spectrogram generation is memory-intensive
MemoryMax=4G
Nice=10

[Install]
WantedBy=multi-user.target
