# systemd/grape-spectrograms.timer
[Unit]
Description=GRAPE Spectrogram Generation Timer (every 10 minutes)
Documentation=https://github.com/mijahauan/grape-recorder

[Timer]
# Run every 10 minutes to keep spectrograms fresh
# Offset by 2 minutes to run after minute boundaries (when new data is available)
OnCalendar=*:02,12,22,32,42,52
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
