# systemd/grape-spectrogram-checker.timer
[Unit]
Description=GRAPE Spectrogram Checker Timer (runs every 30 minutes)
After=network.target

[Timer]
# Run every 30 minutes to keep current day spectrograms fresh
OnCalendar=*:00,30
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
