# systemd/grape-analytics.service
[Unit]
Description=GRAPE Analytics Service - Phase 2 Timing Analysis
Documentation=https://github.com/mijahauan/grape-recorder
After=network-online.target grape-core-recorder.service
Wants=network-online.target
# Analytics depends on core recorder for data
Requires=grape-core-recorder.service

[Service]
Type=simple

# User/Group - same as core recorder
User=grape-recorder
Group=grape-recorder

# Environment from file (set by install.sh)
EnvironmentFile=/etc/grape-recorder/environment

# Working directory
WorkingDirectory=${GRAPE_PROJECT}

# Analytics service launcher script
# This starts all 9 channel analyzers + multi-broadcast fusion
ExecStart=${GRAPE_PROJECT}/scripts/grape-analytics.sh -start ${GRAPE_CONFIG}
ExecStop=${GRAPE_PROJECT}/scripts/grape-analytics.sh -stop

# Don't restart automatically - let individual channel processes handle it
Restart=on-failure
RestartSec=30
StartLimitInterval=300
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grape-analytics

# Resource limits (9 channels + fusion = ~10 processes)
MemoryMax=4G
CPUQuota=200%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only

# Paths we need access to
ReadWritePaths=/var/lib/grape-recorder /var/log/grape-recorder
ReadOnlyPaths=/etc/grape-recorder

[Install]
WantedBy=multi-user.target
