# systemd/grape-radiod-monitor.service
[Unit]
Description=GRAPE Radiod Health Monitor
Documentation=https://github.com/mijahauan/signal-recorder
After=network.target radiod@ac0g-bee1-rx888.service
Wants=radiod@ac0g-bee1-rx888.service

[Service]
Type=simple
User=mjh
Group=mjh
WorkingDirectory=/home/mjh/git/signal-recorder

# Run the radiod health monitor
ExecStart=/home/mjh/git/signal-recorder/venv/bin/python3 \
    /home/mjh/git/signal-recorder/scripts/monitor_radiod_health.py \
    /tmp/grape-test/state/radiod-status.json \
    10

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grape-radiod-monitor

# Process limits
LimitNOFILE=65536
Nice=5

[Install]
WantedBy=multi-user.target
