# systemd/grape-spectrograms.service
[Unit]
Description=GRAPE Spectrogram Generation - Phase 3 Product Update
Documentation=https://github.com/mijahauan/grape-recorder
After=network.target grape-analytics.service

[Service]
Type=oneshot

# User/Group - same as other grape services
User=grape-recorder
Group=grape-recorder

# Environment from file
EnvironmentFile=/etc/grape-recorder/environment

WorkingDirectory=${GRAPE_PROJECT}

# Generate spectrograms for today (rolling 6-hour view, updates every 10 min)
# Uses carrier_spectrogram.py which is the canonical generator
ExecStart=${GRAPE_VENV}/bin/python -m grape_recorder.grape.carrier_spectrogram \
    --data-root ${GRAPE_DATA_ROOT} \
    --all-channels \
    --hours 6

StandardOutput=journal
StandardError=journal
SyslogIdentifier=grape-spectrograms

# Resource limits - spectrogram generation is memory-intensive
MemoryMax=4G
CPUQuota=100%
Nice=10

[Install]
WantedBy=multi-user.target
