# systemd/grape-spectrograms.service
[Unit]
Description=GRAPE Daily Spectrogram Generation
After=network.target

[Service]
Type=oneshot
User=mjh
WorkingDirectory=/home/mjh/git/signal-recorder
Environment="GRAPE_DATA_ROOT=/mnt/grape-data"
ExecStart=/home/mjh/git/signal-recorder/scripts/generate_daily_spectrograms.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grape-spectrograms

# Resource limits
MemoryMax=4G
CPUQuota=50%

[Install]
WantedBy=multi-user.target
