# systemd/grape-daily-upload.service
[Unit]
Description=GRAPE Daily DRF Upload to PSWS
Documentation=https://github.com/wsprdaemon/signal-recorder
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=wsprdaemon
Group=wsprdaemon

# Script location - update PROJECT_ROOT if moved
ExecStart=/home/wsprdaemon/grape-recorder/scripts/daily-drf-upload.sh

# Environment - switch to production when ready
Environment=GRAPE_DATA_ROOT=/tmp/grape-test
Environment=GRAPE_CONFIG=/home/wsprdaemon/grape-recorder/config/grape-config.toml
Environment=GRAPE_VENV=/home/wsprdaemon/grape-recorder/venv
Environment=GRAPE_LOG_DIR=/tmp/grape-test/logs

# Production settings (uncomment when ready):
# Environment=GRAPE_DATA_ROOT=/var/lib/signal-recorder

# Timeouts - allow plenty of time for upload
TimeoutStartSec=3600
TimeoutStopSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grape-daily-upload

[Install]
WantedBy=multi-user.target
