# systemd/signal-recorder.service
[Unit]
Description=GRAPE Signal Recorder
Documentation=https://github.com/yourusername/signal-recorder
After=network-online.target
Wants=network-online.target
# If using radiod, wait for it to start
After=ka9q-radio.service
Wants=ka9q-radio.service

[Service]
Type=simple
User=signal-recorder
Group=signal-recorder

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="MALLOC_ARENA_MAX=2"

# Working directory
WorkingDirectory=/var/lib/signal-recorder

# Configuration file
Environment="CONFIG_FILE=/etc/signal-recorder/config.toml"

# Start command
ExecStart=/usr/local/bin/signal-recorder record --config ${CONFIG_FILE}

# Restart policy
# Restart=on-failure      # Only restart on failure (testing/development)
Restart=always            # Always restart (production - use this!)
RestartSec=30s           # Wait 30 seconds before restarting
StartLimitInterval=300s  # If it fails 5 times in 300s, stop trying
StartLimitBurst=5        # Max 5 restart attempts in the interval

# Resource limits
# Adjust based on number of channels and available memory
MemoryMax=2G
CPUQuota=50%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/signal-recorder /var/log/signal-recorder
ReadOnlyPaths=/etc/signal-recorder

# Capabilities
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_RAW
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_RAW

# Network
# Allow multicast for RTP reception
PrivateNetwork=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=signal-recorder

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
