# systemd/signal-recorder-web.service
[Unit]
Description=GRAPE Signal Recorder Web UI
Documentation=https://github.com/yourusername/signal-recorder
After=network-online.target signal-recorder.service
Wants=network-online.target

[Service]
Type=simple
User=signal-recorder
Group=signal-recorder

# Environment
Environment="NODE_ENV=production"
Environment="PORT=3000"

# Working directory
WorkingDirectory=/opt/signal-recorder/web-ui

# Start command
ExecStart=/usr/bin/node simple-server.js

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitInterval=300s
StartLimitBurst=5

# Resource limits
MemoryMax=512M
CPUQuota=25%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/signal-recorder-web /var/log/signal-recorder
ReadOnlyPaths=/etc/signal-recorder /var/lib/signal-recorder

# Network
# Bind to port 3000
PrivateNetwork=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=signal-recorder-web

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10s

[Install]
WantedBy=multi-user.target
